# deploy/launchd/org.thereisnospoon.seamless.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Seamless launchd service (macOS user LaunchAgent).
  Placeholders __BINARY__ / __CONFIG__ / __LOG__ are substituted by `make install-service`,
  which renders this template to ~/Library/LaunchAgents/ and (re)loads it. Do not hand-edit
  the installed copy; re-run `make install-service` to update.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.thereisnospoon.seamless</string>
  <key>ProgramArguments</key>
  <array>
    <string>__BINARY__</string>
    <string>serve</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SEAMLESS_CONFIG</key>
    <string>__CONFIG__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG__</string>
  <key>StandardErrorPath</key>
  <string>__LOG__</string>
</dict>
</plist>
